# build/image/supervisor/conf.d/iengine.conf
[program:iengine]
command = java -Xmx%(ENV_engineMem)s -jar components/tapdata-agent.jar
environment = app_type="DAAS", tm_port=%(ENV_tm_port)s, backend_url="http://127.0.0.1:%(ENV_tm_port)s/api/", MONGO_URI='%(ENV_MONGO_URI)s', TAPDATA_MONGO_URI='%(ENV_MONGO_URI)s', engineMem=%(ENV_engineMem)s
directory = %(here)s/../../
autostart = false
autoresart = false
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10
stdout_logfile = %(here)s/../logs/iengine.log
stderr_logfile = %(here)s/../logs/iengine.err
stderr_logfile_maxbytes = 50MB
stderr_logfile_backups = 10
