# build/image/supervisor/conf.d/manager.conf
[program:manager]
command = java -Xmx%(ENV_managerMem)s -jar -Dserver.port=%(ENV_tm_port)s -server components/tm.jar --spring.config.additional-location=file:etc/ --logging.config=file:etc/logback.xml --spring.data.mongodb.default.uri=%(ENV_MONGO_URI)s --spring.data.mongodb.obs.uri=%(ENV_MONGO_URI)s --spring.data.mongodb.log.uri=%(ENV_MONGO_URI)s
directory = %(here)s/../../
autostart = false
autoresart = false
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10
stdout_logfile = %(here)s/../logs/manager.log
stderr_logfile = %(here)s/../logs/manager.err
stderr_logfile_maxbytes = 50MB
stderr_logfile_backups = 10
