# build/image/supervisor/conf.d/mongodb.conf
[program:mongodb]
command = mongod --dbpath=/tapdata/data/db/ --replSet=rs0 --wiredTigerCacheSizeGB=%(ENV_dbMem)s --bind_ip_all --logpath=/tapdata/data/logs/mongod.log
autostart = false
autoresart = false
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10
stdout_logfile = %(here)s/../logs/mongodb.log
stderr_logfile = %(here)s/../logs/mongodb.err
stderr_logfile_maxbytes = 50MB
stderr_logfile_backups = 10
